# conf/supervise-mongod.conf
[program:mongod]
command="/usr/bin/mongod"
redirect_stderr=true
autorestart=true
