# conf/supervise-honeypress.conf
[program:honeypress]
command=/opt/honeypress/honeypress.py
stdout_logfile=/opt/honeypress/logs/honey.log
autorestart=true
